# source/services/suricata/etc/init.d/suricata
#!/sbin/openrc-run
# IPCop override for Suricata init script
# Uses /var/run/suricata.pid instead of /var/run/suricata/suricata.pid

name="suricata"
description="Suricata IDS/IPS"

command="/usr/bin/suricata"
pidfile="/var/run/suricata.pid"
# Using -D lets Suricata fork precisely AFTER loading rules, preventing OpenRC false-crashes.
command_args="--pidfile ${pidfile} -D -c /var/ipcop/suricata/suricata.yaml -i eth0"
# Suricata can take a very long time to dump rules to disk during shutdown
retry="TERM/60/KILL/5"

depend() {
	need net
	after firewall
}

start_pre() {
	# Validate configuration before starting
	if [ ! -f "/var/ipcop/suricata/suricata.yaml" ]; then
		eerror "Configuration file not found"
		return 1
	fi
	
	# Create log directory if missing
	mkdir -p /var/log/suricata
	chown root:lighttpd /var/log/suricata
	chmod 0750 /var/log/suricata
}
